# files/system/usr/lib/systemd/system/auto-power-profile.service
[Unit]
Description=Set power profile based on device power supply

[Service]
ExecStart=/usr/bin/auto-power-profile
Restart=on-failure

[Install]
WantedBy=multi-user.target
