# files/system/usr/lib/systemd/system/catcat-flatpak-installer.timer
[Unit]
Description=Timer for installing system flatpaks

[Timer]
OnBootSec=120

[Install]
WantedBy=timers.target
