# files/system/usr/lib/systemd/system/catcat-nightly-reboot.timer
[Unit]
Description=Scheduled Nightly Reboot for CatCat OS

[Timer]
OnCalendar=*-*-* 2:30:00

[Install]
WantedBy=timers.target
