# files/system/usr/lib/systemd/system/catcat-system-setup.service
[Unit]
Description=Configure catcat system

[Service]
Type=simple
ExecStart=/usr/libexec/catcat-system-setup

[Install]
WantedBy=multi-user.target
