# files/system/usr/lib/systemd/system/catcat-system-setup.service
[Unit]
Description=Configure cacat os for current system

[Service]
ExecStart=/usr/libexec/catcat-system-setup

[Install]
WantedBy=default.target
