# files/system/usr/lib/systemd/system/catcat-nightly-reboot.service
[Unit]
Description=Nightly Reboot for CatCat OS

[Service]
Type=oneshot
ExecStart=/usr/libexec/catcat-nightly-reboot
