# files/system/usr/lib/systemd/system/catcat-flatpak-manager.timer
[Unit]
Description=Timer for catcat-flatpak-manager

[Timer]
OnBootSec=120

[Install]
WantedBy=timers.target
