# files/system/usr/lib/systemd/system/catcat-system-maintenance.service
[Unit]
Description=Cleanup and Maintenance of CatCat OS

[Service]
Type=oneshot
ExecStart=/usr/libexec/catcat-system-maintenance

[Install]
WantedBy=multi-user.target
