# files/system/usr/lib/systemd/system/nightly-reboot.timer
[Unit]
Description=Scheduled Nightly Reboot for CatCat OS

[Timer]
OnCalendar=*-*-* 3:30:00

[Install]
WantedBy=timers.target
