# files/system/etc/profile.d/aliases.sh
alias c='clear'
alias ls='ls --color=auto'
alias procs="procs --load-config ${PROCS_CONFIG_FILE}"

