# files/system/etc/profile.d/envar.sh
#export EDITOR="micro"
export MICRO_TRUECOLOR=1
export QT_QPA_PLATFORMTHEME=qt5ct
export QT_QPA_PLATFORM=xcb,wayland
export QT_STYLE_OVERRIDE=kvantum
export GNOME_SHELL_SLOWDOWN_FACTOR=0.6
export JavaScriptCoreUseJIT=0 # https://trac.webkit.org/wiki/EnvironmentVariables
export GJS_DISABLE_JIT=1 # https://gitlab.gnome.org/GNOME/gjs/-/blob/master/doc/Environment.md
