# files/system/etc/profile.d/startup.sh
#!/bin/bash
#/usr/libexec/catcat-user-setup >/dev/null 2>&1
( /usr/bin/bat cache --build >/dev/null 2>&1 & )
